# cpu-on-demand/on-demand-cpupower.service.in
[Unit]
Description=Switch system CPU into powersave mode
ConditionVirtualization=no

[Service]
Type=oneshot
ExecStart=-@SBINDIR@/on-demand-cpupower-switch start

[Install]
WantedBy=multi-user.target
