# cpu-on-demand/on-demand-cpupower.service.in
[Unit]
Description=Switch system into ondemand governer
ConditionVirtualization=no

[Service]
Type=oneshot
ExecStart=-@SBINDIR@/on-demand-cpupower-switch start
ExecStop=-@SBINDIR@/on-demand-cpupower-switch stop

[Install]
WantedBy=multi-user.target
