# cpu-powersave/cpu-powersave.service.in
[Unit]
Description=Switch system CPU into powersave mode
ConditionVirtualization=no

[Service]
Type=simple
ExecStart=-@SBINDIR@/cpu-powersave start
ExecStop=-@SBINDIR@/cpu-powersave stop
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
